# configs/systemd/ollama.service.d/override.conf
# Drop-in override for Ollama systemd service
# INSTALL: sudo mkdir -p /etc/systemd/system/ollama.service.d/
#          sudo cp configs/systemd/ollama.service.d/override.conf /etc/systemd/system/ollama.service.d/override.conf
#          sudo systemctl daemon-reload && sudo systemctl restart ollama
#
# NOTE: Do NOT edit /etc/systemd/system/ollama.service directly.
# Ollama's own installer/updater will overwrite the base unit.
# Use this drop-in overlay to persist settings across Ollama upgrades.
# See ADR-011 in decisions.md for rationale.

[Service]
# Bind to all interfaces so Docker containers (n8n) can reach Ollama via host.docker.internal.
# Required for: The Council, /pm, Weekly News Digest, Linear AI Project Manager.
Environment="OLLAMA_HOST=0.0.0.0"
# Release GPU VRAM after 5 min idle — lets ComfyUI use the GPU between AI requests.
Environment="OLLAMA_KEEP_ALIVE=5m"
# Flash attention for faster inference on RTX 5070 Ti (Ampere+).
Environment="OLLAMA_FLASH_ATTENTION=1"
