# services/openscad-api/openscad-api.service
[Unit]
Description=OpenSCAD Render API — 3D model rendering service for n8n
After=network.target

[Service]
Type=simple
User=biulatech
Group=biulatech
ExecStart=/usr/bin/python3 /home/biulatech/ai-workers-1/services/openscad-api/server.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openscad-api
Environment=PORT=8189

[Install]
WantedBy=multi-user.target
