# configs/systemd/ngrok.service
[Unit]
Description=ngrok tunnel for n8n webhooks
Documentation=https://ngrok.com/docs
After=network-online.target
Wants=network-online.target

[Service]
User=biulatech
ExecStart=/home/biulatech/.local/bin/ngrok http \
  --domain=appendicular-wilson-looser.ngrok-free.dev \
  5678 \
  --log=stdout
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ngrok

[Install]
WantedBy=multi-user.target
